# oo.conf
[program:oo]
command=python3 /home/curator2/oo/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/oo/test.err.log
stdout_logfile=/home/curator2/oo/test.out.log